# userdata/upstart-job.conf
description "forwardPremiumOut dispatcher"

start on stopped rc RUNLEVEL=[2345]

pre-start exec mkdir -p /var/run/lsci2012

exec /usr/bin/python /var/run/lsci2012/fp_dispatch.py
